# etc/supervisord.d/heka.ini
[program:heka]
command=/opt/heka/bin/hekad -config=/etc/heka/hekad.toml
redirect_stderr=true
startsecs=2
exitcodes=0
stdout_logfile=/var/log/supervisor/%(program_name)s.log
